# custom-systemd-units/anki-morning.service
[Unit]
Description=Launch Anki for morning review

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/usr/bin/anki
