# custom-systemd-units/evansync.service
[Unit]
Description=mbsync -a

[Service]
ExecStart=mbsync -a
