# custom-systemd-units/evansync.timer
[Unit]
Description=Periodically run a "cron" script

[Timer]
OnBootSec=15min
OnUnitActiveSec=90min

[Install]
WantedBy=timers.target
