# custom-systemd-units/evil-chin.service
[Unit]
Description=Runs the Evil Chin bot

[Service]
ExecStart=/home/evan/SkyNet/github/evil-chin/run.sh
