# custom-systemd-units/evil-chin.service
[Unit]
Description=Runs the Evil Chin bot

[Service]
ExecStart=/home/evan/Sync/Projects/evil-chin/run.sh
Restart=on-failure
RestartSec=3
