# custom-systemd-units/mosp-2021.service
[Unit]
Description=Runs the MOSP bot suite

[Service]
ExecStart=/home/evan/SkyNet/github/mosp-2021/core/run.sh
