# custom-systemd-units/motioncam.service
[Unit]
Description=Motion camera

[Service]
ExecStart=/home/evan/dotfiles/sentinel/motion.sh
Restart=on-failure
RestartSec=3s
