# custom-systemd-units/motion.service
[Unit]
Description=Motion camera

[Service]
ExecStart=motion -c /home/evan/dotfiles/skynetsys/motion.conf
