# custom-systemd-units/skynet.service
[Unit]
Description=SkyNet server
After=network.target

[Service]
ExecStart=/home/evan/dotfiles/skynetsys/run.sh
