# config/systemd/user/anki-morning.timer
[Unit]
Description=Run Anki at 4am daily

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
