# config/systemd/user/evansync.service
[Unit]
Description=Runs a custom cron script

[Service]
ExecStart=/home/evan/dotfiles/sh-scripts/cron.sh
